# py_battery_status.service
[Unit]
Description=Py Battery Status

[Service]
ExecStart=/home/djagatiya/py_battery/run.sh

[Install]
WantedBy=multi-user.target

